# deploy/narvana-worker.service
[Unit]
Description=Narvana Control Plane Build Worker
After=network.target postgresql.service narvana-api.service
Requires=postgresql.service

[Service]
Type=simple
User=narvana
Group=narvana
WorkingDirectory=/opt/narvana/control-plane
ExecStart=/opt/narvana/control-plane/bin/worker
Restart=always
RestartSec=5

# Environment
EnvironmentFile=/etc/narvana/control-plane.env
Environment=HOME=/opt/narvana
Environment=XDG_RUNTIME_DIR=/run/user/1001

# Security settings relaxed for rootless Podman (requires user namespaces)
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=false
PrivateTmp=false

[Install]
WantedBy=multi-user.target
