# deploy/narvana-web.service
[Unit]
Description=Narvana Control Plane Web UI
After=network.target narvana-api.service

[Service]
Type=simple
User=narvana
Group=narvana
WorkingDirectory=/opt/narvana/control-plane
ExecStart=/opt/narvana/control-plane/bin/web
Restart=always
RestartSec=5

# Environment
EnvironmentFile=/etc/narvana/control-plane.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/log/narvana

[Install]
WantedBy=multi-user.target
